# scim.path
[Unit]
Description=This path will tell you that isf (scim) daemons are ready to launch

[Path]
PathExists=/run/.wm_ready

